# dot_config/systemd/user/rclone-dropbox.service
[Unit]
Description=rclone mount for Dropbox
After=network-online.target
Wants=network-online.target

[Service]
Type=notify
ExecStartPre=-/bin/mkdir -p %h/Dropbox
ExecStart=/usr/bin/rclone mount dropbox: %h/Dropbox \
    --vfs-cache-mode full \
    --allow-other \
    --poll-interval 10s \
    --log-file=%h/.cache/rclone/rclone-dropbox.log \
    --log-level INFO
ExecStop=/bin/fusermount -uz %h/Dropbox

[Install]
WantedBy=default.target